# managed/config/systemd/user/codex-hunt.timer
[Unit]
Description=Run codex hunt every hour

[Timer]
OnBootSec=5min
OnUnitActiveSec=1h
RandomizedDelaySec=120

[Install]
WantedBy=timers.target
