# managed/config/systemd/user/librepods.service
[Unit]
Description=LibrePods (start minimized after GUI startup)
After=graphical-session.target
Wants=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
WorkingDirectory=%h/Documents/Appimages/LibrePods
ExecStart=%h/Documents/Appimages/LibrePods/librepods-x86_64.AppImage --start-minimized
Restart=on-failure
RestartSec=5

[Install]
WantedBy=graphical-session.target
