# managed/config/systemd/user/material-popups.service
[Unit]
Description=Material Popups clipboard overlay
PartOf=graphical-session.target
After=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.local/share/cargo/bin/material-popups
Restart=on-failure
RestartSec=2

[Install]
WantedBy=graphical-session.target
