# managed/config/systemd/user/sync-roommate-sheets.service
[Unit]
Description=Sync Google Sheets to Roommate Finder and notify matches

[Service]
Type=oneshot
WorkingDirectory=/home/magni/Projects/bubin-rooms
ExecStart=/usr/bin/python3 /home/magni/Projects/bubin-rooms/sync-sheet.py
