# .config/systemd/user/cli-proxy-api.service.d/override.conf
[Service]
ExecStart=
ExecStart=%h/.local/share/mise/installs/ubi-router-for-me-cli-proxy-api/latest/cli-proxy-api --config %h/.config/cli-proxy-api/config.yaml
WorkingDirectory=%h/.local/state/cli-proxy-api

ProtectSystem=strict
ProtectHome=tmpfs
PrivateTmp=true
BindPaths=%h/.local/state/cli-proxy-api %h/.config/cli-proxy-api
# mise installs the binary under $HOME; ProtectHome=tmpfs hides it, so re-expose it read-only
BindReadOnlyPaths=%h/.local/share/mise/installs/ubi-router-for-me-cli-proxy-api

NoNewPrivileges=true
CapabilityBoundingSet=
AmbientCapabilities=
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
ProtectClock=true
ProtectHostname=true
ProtectProc=invisible
ProcSubset=pid
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictRealtime=true
RestrictSUIDSGID=true
RestrictNamespaces=true
RemoveIPC=true
UMask=0077

RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallFilter=~@privileged @resources @obsolete
